# c8006/a3/ips.service
[Unit]
Description=IPS
After=sshd.service

[Service]
ExecStart=/usr/local/etc/a3ips/runner.sh

[Install]
WantedBy=multi-user.target
